# systemd/system/snap-snapd-17336.mount
[Unit]
Description=Mount unit for snapd, revision 17336
Before=snapd.service
After=zfs-mount.service

[Mount]
What=/var/lib/snapd/snaps/snapd_17336.snap
Where=/snap/snapd/17336
Type=squashfs
Options=nodev,ro,x-gdu.hide,x-gvfs-hide
LazyUnmount=yes

[Install]
WantedBy=default.target multi-user.target
