# systemd/system/snap-lxd-22248.mount
[Unit]
Description=Mount unit for lxd, revision 22248
Before=snapd.service
After=zfs-mount.service

[Mount]
What=/var/lib/snapd/snaps/lxd_22248.snap
Where=/snap/lxd/22248
Type=squashfs
Options=nodev,ro,x-gdu.hide,x-gvfs-hide
LazyUnmount=yes

[Install]
WantedBy=multi-user.target
